# go_webserver.service
[Unit]
Description=gowebserver
ConditionPathExists=/home/ubuntu/go_webserver/bin/go_webserver
After=network.target
[Service]
Type=simple
User=gowebserver
Group=gowebserver
WorkingDirectory=/home/ubuntu/go_webserver
ExecStart=/home/ubuntu/go_webserver/bin/go_webserver
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=gowebserver
[Install]
WantedBy=multi-user.target
